# screen-log.service
[Unit]
Description=Keep a log of screenshots

[Service]
Type=simple
User=root
Group=root
ExecStart=/home/addison/Documents/repos/screen-log/screenshot.sh

[Install]
WantedBy=multi-user.target

